# scripts/superclock-server.service
[Unit]
Description=SuperClock Express Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nickv2026
WorkingDirectory=/home/nickv2026/SuperClock
EnvironmentFile=-/etc/default/superclock
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
